# pyircx.service
[Unit]
Description=pyIRCX IRC/IRCX Server
Documentation=https://github.com/pyircx/pyircx
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=pyircx
Group=pyircx
WorkingDirectory=/opt/pyircx
ExecStart=/usr/bin/python3 /opt/pyircx/pyircx.py --systemd
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID

# Logging - use journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pyircx

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true

# Allow binding to privileged ports if needed (< 1024)
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# Read-write paths for data
ReadWritePaths=/opt/pyircx/transcripts
ReadWritePaths=/opt/pyircx

# Restart policy
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=3

# Timeouts
TimeoutStartSec=30
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
