# pyircx-webchat.service
[Unit]
Description=pyIRCX WebChat Gateway
Documentation=https://github.com/0x8007000E/pyIRCX
After=network.target pyircx.service
Wants=pyircx.service

[Service]
Type=simple
User=pyircx
Group=pyircx
WorkingDirectory=/opt/pyircx/webchat
EnvironmentFile=-/etc/pyircx/webchat.conf
ExecStart=/usr/bin/python3 /opt/pyircx/webchat/gateway.py --config /etc/pyircx/webchat.conf
Restart=on-failure
RestartSec=5

# Security
NoNewPrivileges=true
PrivateTmp=false
ProtectSystem=false
ProtectHome=true
ReadWritePaths=/opt/pyircx/webchat

[Install]
WantedBy=multi-user.target
